# packaging/systemd/lclva-llama.service
[Unit]
Description=lclva — llama.cpp (LLM backend)
After=network.target

[Service]
Type=simple
ExecStart=%h/.local/opt/llama.cpp/build/bin/llama-server \
  --host 127.0.0.1 \
  --port 8081 \
  --model %h/.local/share/lclva/models/qwen2.5-7b-instruct-q4_k_m.gguf \
  --ctx-size 8192 \
  --n-gpu-layers 999 \
  --threads 4 \
  --metrics
Restart=on-failure
RestartSec=2
TimeoutStartSec=120
# Resource limits — the LLM is the heaviest component.
MemoryMax=10G

[Install]
WantedBy=default.target
